# 95-monitor-hotplug.rules
KERNEL=="card1", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/dan/.Xauthority", RUN+="/home/theruvial/.hotplug_monitor.sh"
